# multi-node/vsphere-powercli/ceph/storage-host3/home-core-data-ceph-osd-ceph-301.mount
[Unit]
Before=kubelet.service
After=format-home-core-data-ceph-osd.service
Requires=format-home-core-data-ceph-osd.service

[Install]
RequiredBy=kubelet.service

[Mount]
What=/dev/sdc1
Where=/home/core/data/ceph/osd/ceph-301
Type=xfs
DirectoryMode=0777